# crontab
# deployed on {{NOW}}

MAILTO=me@josuah.net

#Min	#Hour	#DoM	#Mon	#DoW	#User	#Command

# crontab
0	6	*	*	*	root	calendar -a -A 5
%%! cat /etc/conf/cron/* 2>/dev/null
