# module/monitower/etc/crontab.d/monitower.cron
*	*	*	*	*	root	monitower-run /etc/monitower/queue
