# conf/crontab
# deployed on {{date}}

MAILTO=me@josuah.net

#Min	#Hour	#DoM	#Mon	#DoW	#User	#Command
*/10	*	*	*	*	root	newsyslog
*/5	*	*	*	*	root	rdate -s pool.ntp.org
*/5	*	*	*	*	root	rcctl ls failed
