# conf/crontab
# by {{ path }} at {{ date }}
#Min	#Hour	#DoM	#Mon	#DoW	#User		#Command

*/10	*	*	*	*	root		newsyslog
